# deploy/tk-video-studio.service.example
[Unit]
Description=TK Video Studio (production, port 8000)
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/projects/tk-video-studio
Environment=PATH=/usr/local/bin:/usr/bin:/bin
# --foreground：由 systemd 托管进程，勿在脚本内再 fork 后台
ExecStart=/projects/tk-video-studio/deploy/start-prod.sh --foreground
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
